# dot_config/systemd/user/gdrive-backup.service
[Unit]
Description=Gdrive backup service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/scripts/backup.sh

[Install]
WantedBy=gdrive-backup.timer
